# debian/pymc-repeater.service
[Unit]
Description=PyMC Repeater Daemon
After=network.target

[Service]
Type=simple
User=pymc-repeater
Group=pymc-repeater
WorkingDirectory=/var/lib/pymc_repeater
ExecStart=/usr/bin/python3 -m repeater.main --config /etc/pymc_repeater/config.yaml
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
